# dist/centralbooking.service
# -*- ini -*-

[Unit]
Description=centralbooking service

After=vault.service

[Service]
User=cntrlbook

## provides VAULT_ADDR, VAULT_TOKEN, CONSUL_SERVER_ADDRS
## generated by external source
EnvironmentFile=/var/lib/centralbooking/environment

Environment="LOG_FILE=/var/log/centralbooking/service.log"

ExecStart=/usr/bin/centralbooking
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
